# templates/horizon.ini
[program:horizon]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/opt/php@7.1/bin/php /Users/{{ work_user }}/{{ work_dir }}/{{ dashboard_dir }}/artisan horizon
autostart=true
autorestart=true
user=user
redirect_stderr=true
stdout_logfile=/usr/local/var/log/supervisor_horizon.out.log
